# wgbot.service
[Unit]
Description=nikonVPN Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=wgbot
Group=wgbot
WorkingDirectory=/opt/wgbot
Environment=PATH=/opt/wgbot/venv/bin
ExecStart=/opt/wgbot/venv/bin/python run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wgbot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/wgbot/data /opt/wgbot/logs

[Install]
WantedBy=multi-user.target
